# scripts/crontab
*/20 0-12,21-23 * * * bash /home/dog/rd.sh >/dev/null 2>&1
